# scripts/jarvis-voice.service
[Unit]
Description=Jarvis dashboard + voice loop + live mode + wake word
After=jarvis-vlm.service network-online.target sound.target
# NOT Requires/Wants jarvis-vlm: the dashboard boots STANDALONE so the box comes
# up cool in ECO (the VLM is disabled from auto-start; it's started on the first
# request or "wake up"). jarvis_voice degrades gracefully when the VLM is down
# (_VLM_UP flag → ECO state). After= keeps boot ordering if the VLM is ever up.
Wants=network-online.target

[Service]
Type=simple
User=zip
Group=zip
WorkingDirectory=/home/zip/jarvis-lab
# No VLM-health gate: in ECO the VLM is intentionally down at boot. The dashboard
# serves immediately and reflects ECO; the first request/​wake starts the VLM.
ExecStart=/usr/bin/python3 /home/zip/jarvis-lab/scripts/jarvis_voice.py
Restart=on-failure
RestartSec=3
TimeoutStartSec=90
StandardOutput=append:/home/zip/jarvis-lab/logs/jarvis_voice.log
StandardError=inherit

[Install]
WantedBy=multi-user.target
