# scripts/jarvis-owl.service
[Unit]
Description=Jarvis NanoOWL open-vocab detector sidecar (container)
# Perception mode: OWL and the VLM are MUTUALLY EXCLUSIVE on the 8GB Orin Nano —
# they cannot co-reside (confirmed OOM / VLM crash-loop). Conflicts= enforces the
# swap at the systemd level: starting OWL stops the VLM, and starting the VLM
# stops OWL. (Previously this unit Wants=jarvis-vlm, which pulled the VLM right
# back and OOM'd OWL — the cause of perception-mode detect failures.)
After=docker.service network-online.target
Requires=docker.service
Conflicts=jarvis-vlm.service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=HOME=/home/zip
# Free reclaimable + contiguous memory before the TRT engine grabs the GPU
# (Tegra CMA preflight — runs as root; no VLM-health wait, the VLM is intentionally down).
ExecStartPre=/bin/sh -c 'sync; echo 3 > /proc/sys/vm/drop_caches 2>/dev/null || true; echo 1 > /proc/sys/vm/compact_memory 2>/dev/null || true'
ExecStart=/bin/bash /home/zip/jarvis-lab/scripts/run_owl_sidecar.sh
ExecStop=/usr/bin/docker stop owl-sidecar
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
